# supply-drop-bbs.service
[Unit]
Description=Supply Drop BBS
Documentation=https://github.com/Mesh-America/supply-drop-bbs
After=network.target pymc-companion.service
Wants=network.target
# If using the HAT connection type, pymc-companion must be running first.
# Remove the After= line above if using a USB serial radio instead.

[Service]
Type=simple
User=supply-drop
Group=supply-drop
# dialout gives access to the radio's serial port
SupplementaryGroups=dialout
ExecStart=/usr/local/bin/supply-drop-bbs run --config /etc/supply-drop-bbs/config.toml
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=supply-drop-bbs

# Allow writing to the data directory (set in config.toml data_dir)
ReadWritePaths=/var/lib/supply-drop-bbs
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
